# madsonic-installer-debian/src/DEBIAN/prerm
#! /bin/sh

set -e

# Stop Madsonic service.
if [ -e /etc/init.d/madsonic ]; then
  invoke-rc.d madsonic stop
fi
